# deployment/ufo-galaxy-l4.service
[Unit]
Description=UFO Galaxy L4 Autonomous Intelligence System
After=network.target

[Service]
Type=simple
User=galaxy
WorkingDirectory=/app
ExecStart=/usr/bin/python3 /app/start_l4.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

# 环境变量
Environment="PYTHONUNBUFFERED=1"
Environment="PYTHONPATH=/app"

[Install]
WantedBy=multi-user.target
